# home/.config/autostart/v2rayN-bin.desktop
[Desktop Entry]
Name=v2rayN
Exec=v2rayN
Terminal=false
Type=Application
Icon=v2rayN
Comment=v2rayN Desktop App
Categories=Network
